# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:hugo]
command=/usr/local/bin/hugo server --bind 0.0.0.0 --port 1313 --baseURL http://localhost --watch
directory=/app
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/hugo.err.log
stdout_logfile=/var/log/supervisor/hugo.out.log
environment=HUGO_ENV="production"

[program:admin]
command=/usr/local/bin/node server.js
directory=/app/admin
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/admin.err.log
stdout_logfile=/var/log/supervisor/admin.out.log
environment=NODE_ENV="production",PORT="3001"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
